# crontab
15 8 * * 5 cd /app && python -m bg_tasks.update_performance_weekly > /app-logs/update_performance_weekly.log
5 */8 * * * cd /app && python -m bg_tasks.update_delta_neutral_vault_performance_daily > /app-logs/update_delta_neutral_vault_performance_daily.log
# 0 12 * * * cd /app && python -m bg_tasks.update_lido_reward > /app-logs/update_lido_reward.log